# SWD & Provisioning Scripts/SDCache Purge/ldcron-sdclean.xml
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task"><Triggers><TimeTrigger id="1"><Enabled>true</Enabled><StartBoundary>2016-09-28T06:10:44-0600</StartBoundary><Repetition><Interval>PT604800S</Interval></Repetition></TimeTrigger></Triggers><Queues><Queue>sdclean</Queue></Queues><Actions><Exec><Command>find /Library/Application\ Support/LANDesk/sdcache/* -mtime +10 -exec rm -rf {} \;</Command></Exec></Actions></Task>